# motion-snapshot.timer
[Unit]
Description=Run motion snapshot capture/upload at specific times

[Timer]
# Template defaults.
# During installation, install.sh rewrites OnCalendar entries based on
# interactive HH:MM schedule times you provide.
OnCalendar=*-*-* 06:00:00
OnCalendar=*-*-* 14:00:00
OnCalendar=*-*-* 22:00:00
Persistent=true
Unit=motion-snapshot.service

[Install]
WantedBy=timers.target